# debian/pmg-hourly.timer
[Unit]
Description=Hourly Proxmox Mail Gateway activities

[Timer]
OnCalendar=hourly
Persistent=true     

[Install]
WantedBy=timers.target