# debian/pmg-cron-hourly.service
[Unit]
Description=Proxmox Mail Gateway Hourly Tasks

[Service]
Type=oneshot
ExecStart=/usr/lib/pmg/bin/pmg-cron-hourly
