# debian/pmg-daily.service
[Unit]
Description=Daily Proxmox Mail Gateway activities

[Service]
Type=oneshot
ExecStart=/usr/lib/pmg/bin/pmg-daily
