# debian/pmg-cron-hourly.timer
[Unit]
Description=Run pmg-cron-hourly.service hourly

[Timer]
OnCalendar=hourly
Persistent=true     

[Install]
WantedBy=timers.target