# debian/pmgreport.timer
[Unit]
Description=Send Daily System Report Mail

[Timer]
OnCalendar=00:01
Persistent=true

[Install]
WantedBy=timers.target