# debian/pmgspamreport.timer
[Unit]
Description=Send Daily Spam Report Mails

[Timer]
OnCalendar=00:05
Persistent=true  

[Install]
WantedBy=timers.target